# com.coach.bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.coach.bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/I048171/Claude/coach/bot.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/I048171/Claude/coach</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/I048171/Claude/coach/logs/bot.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/I048171/Claude/coach/logs/bot.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/I048171/Library/Python/3.9/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/I048171/Claude/coach</string>
    </dict>
</dict>
</plist>
